# xmpp-proxy-stack/supervisord.conf
[supervisord]
nodaemon=true
logfile=/logs/supervisor.log
pidfile=/var/run/supervisord.pid
user=root

[program:xmpp-proxy]
command=/usr/local/bin/xmpp-proxy /etc/xmpp-proxy/xmpp-proxy.toml
autostart=true
autorestart=true
stdout_logfile=/logs/xmpp-proxy.log
stderr_logfile=/logs/xmpp-proxy.log
redirect_stderr=true
user=root
priority=10

[program:fail2ban-rs]
command=/usr/local/bin/fail2ban-rs --config /etc/fail2ban-rs/config.toml
autostart=true
autorestart=true
stdout_logfile=/logs/fail2ban-rs.log
stderr_logfile=/logs/fail2ban-rs.log
redirect_stderr=true
user=root
priority=20
